# systemd/system/nominatim.service.d/override.conf
[Unit]
After=postgresql.service
Requires=postgresql.service
